# templates/upstart.conf.erb
start on runlevel [2345]
stop on runlevel [!2345]

chdir <%= scope.lookupvar("solr::params::jetty_dir") %>

script
	sudo -n -u <%= scope.lookupvar("solr::params::user") %> java <%= @java_opts %> -jar <%= scope.lookupvar("solr::params::var_dir") %>/jetty/start.jar <%= scope.lookupvar("solr::params::conf_dir") %>/jetty.xml | logger -t solr
end script